# reka-deploy/roles/reka/templates/reka.upstart.conf.j2
description "Reka service"
author "Nick Sellen"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env CONTAINER=production
chdir /var/apps/reka
exec bin/reka-server >> reka.log 2>&1